# os/macos/crontab/crontab
0 2 * * * /usr/local/bin/brew update >> ~/logs/cron/brew/update.log 2>&1
